# install/systemd/ollama@.service
[Unit]
Description=Ollama LLM Service
After=network.target

[Service]
Type=simple
User=%i
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="OLLAMA_HOST=127.0.0.1:11434"
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
